# certbot.crontab
0 0 * * *	root 	cd /home/nginx/nginx-proxy ; docker-compose run --service-ports certbot 0</dev/null && docker-compose kill -s 'SIGHUP' nginx
